# roles/web/files/etc/cron.d/sessionclean
12,42 * * * *   root    /usr/local/sbin/sessionclean
